# debian/onak.socket
[Socket]
ListenStream=/run/keyd.sock
SocketUser=onak
SocketMode=0660

[Install]
WantedBy=sockets.target
